# dicod-docker.service
[Unit]
Description=Dockerized GNU Dico DICT server

[Service]
ExecStart=/home/aporia/scripts/dicod-docker

[Install]
WantedBy=multi-user.target

