# muse-parse.service
[Unit]
Description=Parse muse reading logs on boot

[Service]
ExecStart=/home/aporia/.local/bin/muse parse

[Install]
WantedBy=multi-user.target
